# scripts/upstart/bdemo.conf
description "Beego-demo service"

start on runlevel [2345]
stop on runlevel [06]

limit nofile 65535 65535
chdir /var/beego-demo

exec ./beego-demo >>./logs/console.log 2>&1
